# templates/dev.pf.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
<dict>
    <key>Label</key>
    <string><%= scope.lookupvar('servicename') %></string>
    <key>Program</key>
    <string>/sbin/pfctl</string>
    <key>ProgramArguments</key>
    <array>
        <string>/sbin/pfctl</string>
        <string>-e</string>
        <string>-f</string>
        <string><%= scope.lookupvar('configdir') %>/pf.conf</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>ServiceDescription</key>
    <string>Boxen's FreeBSD Packet Filter (pf) daemon</string>
    <key>StandardErrorPath</key>
    <string><%= scope.lookupvar('boxen::config::logdir') %>/pf.log</string>
    <key>StandardOutPath</key>
    <string><%= scope.lookupvar('boxen::config::logdir') %>/pf.log</string>
</dict>
</plist>
